# creatornewsdesk-dashboard.service
[Unit]
Description=Creator Newsdesk Dashboard
After=network.target

[Service]
Type=simple
User=mmurphy
WorkingDirectory=/home/mmurphy/homelab/projects/creatornewsdesk
ExecStart=/usr/bin/python3 /home/mmurphy/homelab/projects/creatornewsdesk/dashboard_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
